# transmission-nas/rootfs/etc/services.d/transmission/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Transmission
# ==============================================================================
bashio::log.info "Starting Transmission..."

# Mount CIFS Share if configured and if Protection Mode is active
if bashio::config.has_value 'networkdisks'; then
    MOREDISKS=$(bashio::config 'networkdisks')
    CIFS_USERNAME=$(bashio::config 'cifsusername')
    CIFS_PASSWORD=$(bashio::config 'cifspassword')
    if bashio::config.has_value 'cifsversion'; then
        CIFS_VERSION_ARG=",vers=$(bashio::config 'cifsversion')"
    else
        CIFS_VERSION_ARG=""
    fi
    bashio::log.info "Network Disks mounting.. ${MOREDISKS}" && \
    for disk in $MOREDISKS 
    do
        bashio::log.info "Mount ${disk}"
        mkdir -p /$disk && \
            mount -t cifs -o username=$CIFS_USERNAME,password=$CIFS_PASSWORD$CIFS_VERSION_ARG $disk /$disk && \
            bashio::log.info "Success!"   
    done || \
    bashio::log.warning "Protection mode is ON. Unable to mount external drivers!"
fi

if bashio::config.true 'openvpn_enabled'; then
    exec /usr/sbin/openvpn --config /etc/openvpn/config.ovpn --script-security 2 --up /etc/openvpn/up.sh --down /etc/openvpn/down.sh
else
    exec /usr/bin/transmission-daemon --foreground --config-dir /data/transmission
fi
